# zfts_html/etc/systemd/zcompd-dops-p2.service
# /etc/systemd/system/zcompd.service
[Unit]
Description=Zfts transfer system dops-p2
After=network.target
Wants=network-online.target

[Service]
Type=simple
UMask=0022
#User=admin
#Group=admin
ExecStart=/opt/prgsnap/bin/zcompd -c /etc/zfts.configs/zfts_zcomp_gnd_dops_p2.ini
ExecStop=/bin/kill -TERM $MAINPID
PIDFile=/run/zcompd.pid
WorkingDirectory=/opt/prgsnap
Restart=on-failure
RestartSec=5
# Optional security hardening
NoNewPrivileges=true
PrivateTmp=true

SyslogIdentifier=zcompd-dops-p2

[Install]
WantedBy=default.target
