# rules.d/98-v-274878.rules
#V-274878
-w /etc/cron.d -p wa -k cronjobs
-w /var/spool/cron -p wa -k cronjobs
